# change_wallpaper/change_wallpaper.service
[Unit]
Description=Change the gnome wallpapper

[Service]
ExecStart=SCRIPT_PATH

[Install]
WantedBy=multi-user.target
 
