# turn_off_low_battery/sec_poweroff.service
[Unit]
Description=Turn off on low battery

[Service]
ExecStart=SCRIPT_PATH

[Install]
WantedBy=multi-user.target
 
